# docker-prune.timer
[Unit]
Description=Runs daily Docker system prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
